# common/roles/activate_services/files/prune-docker-contents.service
[Unit]
Description=Prune unneeded docker contents (images, networks, etc)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker system prune -f
