# sngsk.info/roles/install_mackerel_conf_for_nspawn_container/files/set-working-on-mackerel.service
[Unit]
Description=Set host state on Mackerel to 'working'

[Service]
Type=oneshot
ExecStart=/usr/bin/mkr update --status working

[Install]
WantedBy=default.target
