# sngsk.info/roles/install_mackerel_conf_for_sngsk.info/files/set-maintenance-on-mackerel.service
[Unit]
Description=Set host state on Mackerel to 'maintenance'

[Service]
Type=oneshot
ExecStart=/usr/bin/mkr update --status maintenance

[Install]
WantedBy=default.target
