# sngsk.info/roles/activate_services/files/update_container.timer
[Unit]
Description=update arch linux nspawn container (daily)

[Timer]
OnCalendar=daily
RandomizedDelaySec=1day
Persistent=true

[Install]
WantedBy=timers.target
