# sngsk.info/roles/activate_services/files/update_container.service
[Unit]
Description=update arch linux nspawn container

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_arch_nspawn_container
