# sngsk.info/roles/activate_services/files/update-container.service
[Unit]
Description=update arch linux nspawn container

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-arch-nspawn-container
