# sngsk.info/roles/activate_services/files/wake-up-container.timer
[Unit]
Description=wake up container installed in the host (per 5mins)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
