# sngsk.info/roles/activate_services/files/wake-up-container.timer
[Unit]
Description=wake up nginx

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
