# sngsk.info/roles/activate_services/files/wake_up_nginx.service
[Unit]
Description=wake up nginx

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wake_up_nginx
