# sngsk.info/roles/activate_services/files/wake-up-nginx.timer
[Unit]
Description=wake up nginx (per 2mins)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
