# sngsk.info/roles/activate_services/files/wake-up-nginx.service
[Unit]
Description=Making nginx woken up

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wake-up-nginx
