# sngsk.info/roles/activate_services/files/wake-up-nginx.timer
[Unit]
Description=Making nginx waked up

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
