# nspawn-host/roles/activate_services/files/update-container.service
[Unit]
Description=Updating Arch Linux nspawn containers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-arch-nspawn-container
