# nspawn-host/roles/activate_services/files/wake-up-container.service
[Unit]
Description=Makeing containers installed in the host waked up

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wake-up-container
