# nspawn-host/roles/activate_services/files/wake-up-container.service
[Unit]
Description=wake up container installed in the host

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wake-up-container
